# deployment/systemd/projectmeats-socket.service
[Unit]
Description=ProjectMeats Django Application
Documentation=https://github.com/Vacilator/ProjectMeats
After=network.target postgresql.service
Wants=postgresql.service
Requires=projectmeats.socket

[Service]
Type=notify
User=projectmeats
Group=www-data
RuntimeDirectory=projectmeats
RuntimeDirectoryMode=0755
WorkingDirectory=/opt/projectmeats/backend
Environment=PATH=/opt/projectmeats/venv/bin
Environment=DJANGO_SETTINGS_MODULE=apps.settings.production
Environment=PYTHONPATH=/opt/projectmeats/backend
EnvironmentFile=-/etc/projectmeats/projectmeats.env
EnvironmentFile=-/opt/projectmeats/.env.production
EnvironmentFile=-/opt/projectmeats/backend/.env
ExecStartPre=/opt/projectmeats/deployment/scripts/pre_start_service.sh
ExecStart=/opt/projectmeats/venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --preload \
    --bind unix:/run/projectmeats.sock \
    --access-logfile /var/log/projectmeats/access.log \
    --error-logfile /var/log/projectmeats/error.log \
    --log-level info \
    --timeout 120 \
    --graceful-timeout 30 \
    projectmeats.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
ExecStopPost=/bin/sh -c 'if [ -w /var/log/projectmeats/ ]; then journalctl -u projectmeats.service -n 50 > /tmp/post_failure.log && mv /tmp/post_failure.log /var/log/projectmeats/post_failure.log; fi'
KillMode=mixed
TimeoutStopSec=5
SuccessExitStatus=0 1 2
PrivateTmp=true
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/projectmeats/backend/media /var/log/projectmeats /var/run/projectmeats /run
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target